# embedded-files/etc/uci-defaults/20-freifunk-berlin-zram
#!/bin/sh

uci set system.@system[0].zram_comp_algo=lzo
uci commit system

