# embedded-files/etc/uci-defaults/freifunk-berlin-restore-service-state
#!/bin/sh

SERVICE_STATUS=/etc/backup/service_status.txt

[ -e ${SERVICE_STATUS} ] || exit 0

. ${SERVICE_STATUS}
