# embedded-files/etc/uci-defaults/wipe_glinetconfigs
#!/bin/sh /etc/rc.common

boot() {
  echo "previous GLinet sysstem detected ... flushing flash ..."
  [ -f /etc/config/glconfig ] && firstboot -y -r
}

boot
